# wiroc-rtc-writeback.service
# Persist a known-good (NTP-synced) system clock into the battery-backed pcf8563 RTC
[Unit]
Description=WiRoc: persist system time to the pcf8563 RTC (only when NTP-synced)
After=chrony.service time-sync.target

[Service]
Type=oneshot
ExecStart=/bin/sh /home/chip/WiRoc-StartupScripts/wiroc-rtc-writeback.sh
